# www-apps/ablog/files/ablog.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

depend() {
        need net
}

command=/usr/bin/ablog
command_args=""
command_background="true"
pidfile="/run/${RC_SVCNAME}.pid"
name="A blog Daemon"
description="A blog is yey another blog app"

export TMPDIR="${TMPDIR:-/tmp}"
